# upstart/gissues.conf
# this is from http://kevin.vanzonneveld.net/techblog/article/run_nodejs_as_a_service_on_ubuntu_karmic/

description "gissues server"
author      "Tomasz Janczuk <tomasz@janczuk.org>"
 
# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"
 
    exec /usr/local/bin/node /home/tjanczuk/gissues/src/server.js >> /var/log/gissues.log 2>&1
end script
 
post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script